# data/unity-greeter-session-broadcast-session.conf
description "Unity Greeter Session Broadcast Session Component"
author "Ted Gould <ted@canonical.com>"

task

emits application-start

#TODO: Check GSettings Key
start on dbus-system SIGNAL="StartApplication" INTERFACE="com.canonical.Unity.Greeter.Broadcast" PATH="/com/canonical/Unity/Greeter/Broadcast" ARG0=${USER}

pre-start script
	OWNER=`gdbus call --system --dest org.freedesktop.DBus --object-path / --method org.freedesktop.DBus.GetNameOwner com.canonical.Unity.Greeter.Broadcast | cut -d "'" -f 2`

	if [ -z $OWNER ] ; then
		stop
		exit 1
	fi
	if [ $OWNER != $SENDER ] ; then
		stop
		exit 1
	fi
end script


exec initctl emit application-start APP_ID=${ARG1}
