# data/unity-greeter-session-broadcast-session.conf
description "Unity Greeter Session Broadcast Session Component"
author "Ted Gould <ted@canonical.com>"

task

emits application-start

#TODO: Check GSettings Key
start on dbus-system SIGNAL="StartApplication" INTERFACE="com.canonical.Unity.Greeter.Broadcast" SENDER="com.canonical.Unity.Greeter.Broadcast" PATH="/com/canonical/Unity/Greeter/Broadcast" ARG0=${USER}

exec initctl emit application-start APP_ID=${ARG1}
